# dist/dist-packages/linux/openziti-controller/ziti-controller.service
[Unit]
Description=OpenZiti Controller
After=network-online.target

[Service]

# this env file configures the service, including whether or not to perform bootstrapping
EnvironmentFile=/opt/openziti/etc/controller/service.env

#
## misc
#

# manage the user and permissions for the service automatically
DynamicUser=yes
# used by bootstrap.bash to look up /run/credentials/$UNIT_NAME/$CREDENTIAL_NAME
Environment=UNIT_NAME=ziti-controller.service
# create a new config file relative to working directory unless it exists
# absolute path where service will be run
ExecStartPre=/opt/openziti/etc/controller/entrypoint.bash check config.yml
ExecStart=/opt/openziti/bin/ziti controller run config.yml ${ZITI_ARGS}
LimitNOFILE=65535
Restart=always
RestartSec=3
# relative to /var/lib
StateDirectory=ziti-controller
#:/opt/openziti/etc/controller/state
# "ziti controller run" is the main process managed by this service and replaces entrypoint.bash
Type=simple
UMask=0007
WorkingDirectory=/var/lib/ziti-controller

[Install]
WantedBy=multi-user.target
