# dist/dist-packages/linux/openziti-controller/ziti-controller.service
[Unit]
Description=OpenZiti Controller
After=network-online.target

[Service]

# you must set the admin password to bootstrap the database; the bootstrapping script will initialize the default admin
# username and password by loading the password from a file or env var; for security, set permissions to allow read only
# by root only and remove the file after the bootstrapping that occurs on the first run or set ZITI_PWD in
# /opt/openziti/etc/controller/bootstrap.env
#
LoadCredential=ZITI_PWD:/opt/openziti/etc/controller/.pwd

#
## Optional Permissions
#

# allow binding low ports, e.g., 443/tcp; NOTE: use TLS passthrough if fronting with a reverse proxy, i.e., "raw" TCP
# proxy
# AmbientCapabilities=CAP_NET_BIND_SERVICE

#
##
#

# this env file configures the service, including whether or not to perform bootstrapping
EnvironmentFile=/opt/openziti/etc/controller/service.env
# this answer file provides input values for bootstrapping during first startup
EnvironmentFile=/opt/openziti/etc/controller/bootstrap.env
# manage the user and permissions for the service automatically

DynamicUser=yes
# used by bootstrap.bash to look up /run/credentials/$UNIT_NAME/$CREDENTIAL_NAME
Environment=UNIT_NAME=ziti-controller.service
# create a new config file relative to working directory unless it exists
# absolute path where service will be run
ExecStart=/opt/openziti/etc/controller/entrypoint.bash run config.yml
LimitNOFILE=65535
Restart=always
RestartSec=3
# relative to /var/lib
StateDirectory=ziti-controller
#:/opt/openziti/etc/controller/state
# "ziti controller run" is the main process managed by this service and replaces entrypoint.bash
Type=simple
UMask=0007
WorkingDirectory=/var/lib/ziti-controller

[Install]
WantedBy=multi-user.target
